# qemu/files/runqemu_start@.service
[Unit]
# For use with .path
Description=Start VM %I

[Service]
Type=oneshot
ExecStart=/bin/systemctl start runqemu@%I
