# qemu/files/runqemu_stop@.service
[Unit]
# For use with .path
Description=Stop VM %I

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop runqemu@%I
