# monitorme/files/balance.service
[Unit]
Description=Submit Bitcoin balance to statsd

[Service]
Type=simple
ExecStart=/usr/local/bin/balance

[Install]
WantedBy=multi-user.target
