# monitorme/files/balance.timer
[Unit]
Description=Check balance every minute

[Timer]
OnCalendar=hourly
Unit=balance.service

[Install]
WantedBy=multi-user.target
