# sslexpirystatsd/files/sslexpirystatsd.timer
[Unit]
Description=sslexpirystatsd minutely
[Timer]
OnCalendar=minutely
[Install]
WantedBy=multi-user.target
