# sslexpirystatsd/files/sslexpirystatsd@.timer
[Unit]
Description=sslexpirystatsd %I minutely
[Timer]
OnCalendar=minutely
[Install]
WantedBy=multi-user.target
