# vmmanagement/files/vmmanagement_run_topup.path
[Unit]
Description=Monitors for topup-worthy files

[Path]
Unit=vmmanagement_run_topup.service
DirectoryNotEmpty=/var/tmp/vmmanagement_topup

[Install]
WantedBy=multi-user.target
